# assets/master/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:salt-master]
command=/usr/bin/salt-master --log-file=/var/log/salt/salt-master.log --log-file-level=debug
